# docker/root/etc/services.d/wol_api/run
#!/usr/bin/with-contenv bash

exec wol_api --port=8080
